# rc.d/arkictrl.in
#!/bin/sh
#
# PROVIDE: arkictrl
# REQUIRE: NETWORKING LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="arkictrl"
rcvar="arkictrl_enable"

load_rc_config "$name"

: ${arkictrl_enable:="NO"}
: ${arkictrl_user:="@ARKICTRL_SERVICE_USER@"}
: ${arkictrl_env:="DISPLAY=:0"}

command="/usr/sbin/daemon"

if [ "$(id -u)" -eq 0 ]; then
	# Started as root: daemon(8) drops to ${arkictrl_user} before writing
	# the -p pidfile, so the pidfile must live in a directory that user can
	# write. /var/run is root-only, so use a per-service subdirectory that
	# start_precmd creates and chowns to ${arkictrl_user}.
	: ${arkictrl_pidfile:="/var/run/${name}/${name}.pid"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -p ${pidfile} -u ${arkictrl_user} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
else
	# Started as an unprivileged user (e.g. service arkictrl onestart):
	# /var/run is not writable and daemon -u needs root, so keep the
	# pidfile in a user-writable location and run as the current user.
	: ${arkictrl_pidfile:="${TMPDIR:-/tmp}/${name}-$(id -un).pid"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -p ${pidfile} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
fi

start_precmd="arkictrl_precmd"

arkictrl_precmd()
{
	# Ensure the pidfile directory exists and is writable by the user the
	# daemon will run as (it writes the pidfile after dropping privileges).
	piddir="${pidfile%/*}"
	if [ "$(id -u)" -eq 0 ]; then
		/usr/bin/install -d -o "${arkictrl_user}" -m 0755 "${piddir}"
	else
		/usr/bin/install -d -m 0755 "${piddir}"
	fi
}

run_rc_command "$1"
