# rc.d/arkictrl.in
#!/bin/sh
#
# PROVIDE: arkictrl
# REQUIRE: NETWORKING LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="arkictrl"
rcvar="arkictrl_enable"

load_rc_config "$name"

: ${arkictrl_enable:="NO"}
: ${arkictrl_user:="@ARKICTRL_SERVICE_USER@"}
: ${arkictrl_env:="DISPLAY=:0"}

# rc runs with a minimal PATH that omits /usr/local/{s,}bin, where python3
# and the tools app.py shells out to (pgrep, kdeconnect-cli, ...) live. Make
# them reachable, and apply arkictrl_env (e.g. DISPLAY) to the daemon's
# environment, which is inherited by the spawned process.
export PATH="/usr/local/sbin:/usr/local/bin:${PATH}"
for _kv in ${arkictrl_env}; do
	export "${_kv}"
done

command="/usr/sbin/daemon"

if [ "$(id -u)" -eq 0 ]; then
	# Started as root: daemon(8) drops to ${arkictrl_user} before writing
	# the -p pidfile and -o logfile, so both must live in a directory that
	# user can write. /var/run is root-only, so use a per-service
	# subdirectory that start_precmd creates and chowns to ${arkictrl_user}.
	: ${arkictrl_pidfile:="/var/run/${name}/${name}.pid"}
	: ${arkictrl_logfile:="/var/run/${name}/${name}.log"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -o ${arkictrl_logfile} -p ${pidfile} -u ${arkictrl_user} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
else
	# Started as an unprivileged user (e.g. service arkictrl onestart):
	# /var/run is not writable and daemon -u needs root, so keep the
	# pidfile/logfile in a user-writable location and run as the current user.
	: ${arkictrl_pidfile:="${TMPDIR:-/tmp}/${name}-$(id -un).pid"}
	: ${arkictrl_logfile:="${TMPDIR:-/tmp}/${name}-$(id -un).log"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -o ${arkictrl_logfile} -p ${pidfile} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
fi

start_precmd="arkictrl_precmd"

arkictrl_precmd()
{
	# Ensure the pidfile/logfile directories exist and are writable by the
	# user the daemon will run as (it writes them after dropping privileges).
	if [ "$(id -u)" -eq 0 ]; then
		for d in "${arkictrl_pidfile%/*}" "${arkictrl_logfile%/*}"; do
			/usr/bin/install -d -o "${arkictrl_user}" -m 0755 "${d}"
		done
	else
		for d in "${arkictrl_pidfile%/*}" "${arkictrl_logfile%/*}"; do
			/usr/bin/install -d -m 0755 "${d}"
		done
	fi
}

run_rc_command "$1"
