# rc.d/arkictrl.in
#!/bin/sh
#
# PROVIDE: arkictrl
# REQUIRE: NETWORKING LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="arkictrl"
rcvar="arkictrl_enable"

load_rc_config "$name"

: ${arkictrl_enable:="NO"}
: ${arkictrl_user:="@ARKICTRL_SERVICE_USER@"}
: ${arkictrl_env:="DISPLAY=:0"}

command="/usr/sbin/daemon"

if [ "$(id -u)" -eq 0 ]; then
	# Started as root: write the pidfile to /var/run and drop privileges
	# to ${arkictrl_user} (daemon -u requires superuser).
	: ${arkictrl_pidfile:="/var/run/${name}.pid"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -p ${pidfile} -u ${arkictrl_user} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
else
	# Started as an unprivileged user (e.g. service arkictrl onestart):
	# /var/run is not writable and daemon -u needs root, so keep the
	# pidfile in a user-writable location and run as the current user.
	: ${arkictrl_pidfile:="${TMPDIR:-/tmp}/${name}-$(id -un).pid"}
	pidfile="${arkictrl_pidfile}"
	command_args="-f -p ${pidfile} @CMAKE_INSTALL_PREFIX@/bin/arkictrl"
fi

run_rc_command "$1"
