# systemd/arkictrl-user.service.in
[Unit]
Description=Arki Control – web remote-control interface
After=network.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/arkictrl
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
