# systemd/omdrcctrl-user.service.in
[Unit]
Description=OMDRC Control – web remote-control interface
After=network.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/omdrcctrl
Restart=on-failure
RestartSec=5
KillMode=process

[Install]
WantedBy=default.target
